# runit/service/openvpn/run
#!/bin/bash

exec /sbin/openvpn --config /etc/openvpn/clientcert.ovpn
